# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:redis]
command=redis-server
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0



[program:governance]
command=uvicorn arhiax_dx.main:app --host 0.0.0.0 --port 8088
directory=/app/back
environment=PYTHONUNBUFFERED=1,ARHIAX_DX_SPECS_PATH="/app/back/specs",ARHIAX_DX_LEDGER_PATH="/app/back/var/evidence-ledger.jsonl",HOST="0.0.0.0",PORT="8088"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:api]
command=uvicorn api.main:app --host 0.0.0.0 --port 8000
directory=/app/back-api/src
environment=PYTHONUNBUFFERED=1,DATABASE_URL="%(ENV_DATABASE_URL)s",REDIS_URL="redis://localhost:6379/0",GOVERNANCE_API_URL="http://localhost:8088",DXPRO_URL="http://localhost:8310",SECRET_KEY="arhiax-dx-dev-secret-change-in-prod",ALGORITHM="HS256",ACCESS_TOKEN_EXPIRE_MINUTES="480",APP_URL="%(ENV_APP_URL)s",SPECS_PATH="/app/back/specs",LEDGER_PATH="/app/back/var/evidence-ledger.jsonl",GEMINI_API_KEY="%(ENV_GEMINI_API_KEY)s"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=python -m api.worker
directory=/app/back-api/src
environment=PYTHONUNBUFFERED=1,DATABASE_URL="%(ENV_DATABASE_URL)s",REDIS_URL="redis://localhost:6379/0",GOVERNANCE_API_URL="http://localhost:8088",SECRET_KEY="arhiax-dx-dev-secret-change-in-prod",SPECS_PATH="/app/back/specs",LEDGER_PATH="/app/back/var/evidence-ledger.jsonl"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:dxpro]
command=python -m dxpro_runtime.server
directory=/app/ARHIAX-Dx-Pro
environment=PYTHONUNBUFFERED=1,DXPRO_HOST="0.0.0.0",DXPRO_PORT="8310",DXPRO_ENV="development",DXPRO_LEDGER_PATH="/app/ARHIAX-Dx-Pro/data/evidence.jsonl",DXPRO_CASE_STORE_ROOT="/app/ARHIAX-Dx-Pro/data/cases",DXPRO_EXPORT_ROOT="/app/ARHIAX-Dx-Pro/data/exports",DXPRO_POLICY_BUNDLE_PATH="/app/ARHIAX-Dx-Pro/policy-bundle-pmel-v1.0.0"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:frontend]
command=node server.js
directory=/app/front/.next/standalone
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0",NEXT_PUBLIC_API_URL="http://localhost:8000"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
